# roles/horde_stats_exporter/templates/horde-downsample.service.j2
[Unit]
Description=AI Horde Metrics Downsampling
After=network.target prometheus.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ exporter_user }}
Group={{ exporter_user }}
WorkingDirectory={{ exporter_install_dir }}
ExecStart={{ exporter_install_dir }}/.venv/bin/ai-horde-downsample \
  --source-tenant {{ exporter_downsample_source_tenant }} \
  --target-tenant {{ exporter_downsample_target_tenant }} \
  --mimir-url {{ exporter_downsample_mimir_url }} \
  --resolution {{ exporter_downsample_resolution }}
StandardOutput=append:{{ exporter_log_file }}
StandardError=append:{{ exporter_log_file }}

[Install]
WantedBy=multi-user.target
