# roles/horde_stats_exporter/templates/horde-downsample.timer.j2
[Unit]
Description=Daily AI Horde Metrics Downsampling Timer
Requires=horde-downsample.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
