# roles/horde_stats_exporter/templates/horde-exporter.service.j2
[Unit]
Description=AI Horde Prometheus Exporter
After=network.target
Wants=network-online.target

[Service]
Type=simple
User={{ exporter_user }}
Group={{ exporter_user }}
WorkingDirectory={{ exporter_install_dir }}
ExecStart={{ exporter_install_dir }}/venv/bin/horde-exporter --config {{ exporter_install_dir }}/exporter_config.yaml
Restart=always
RestartSec=10
StandardOutput=append:{{ exporter_log_file }}
StandardError=append:{{ exporter_log_file }}

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ exporter_log_file }}

[Install]
WantedBy=multi-user.target
