# roles/horde_monitoring/templates/prometheus-backup.timer.j2
[Unit]
Description=Prometheus Daily Backup Timer
Requires=prometheus-backup.service

[Timer]
# Run daily at {{ prometheus_backup_schedule_hour | default('02') }}:{{ prometheus_backup_schedule_minute | default('00') }}
OnCalendar={{ prometheus_backup_schedule | default('*-*-* 02:00:00') }}
Persistent=true
RandomizedDelaySec={{ prometheus_backup_random_delay | default('300') }}

[Install]
WantedBy=timers.target
