# roles/horde_monitoring/templates/prometheus-backup.service.j2
[Unit]
Description=Prometheus Backup Service
After=prometheus.service
Requires=prometheus.service

[Service]
Type=oneshot
User=prometheus
Group=prometheus
ExecStart=/usr/local/bin/prometheus_backup.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=prometheus-backup

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ prometheus_backup_dir | default('/var/backups/prometheus') }} {{ prometheus_backup_data_dir | default('/var/lib/prometheus') }} /var/log

[Install]
WantedBy=multi-user.target
