# roles/horde_monitoring/templates/prometheus.service.j2
[Unit]
Description=Prometheus Time Series Database
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=prometheus
Group=prometheus
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=/usr/local/bin/prometheus \
  --config.file=/etc/prometheus/prometheus.yml \
  --storage.tsdb.path=/var/lib/prometheus/ \
{% if prometheus_retention_time | default('0') != '0' %}
  --storage.tsdb.retention.time={{ prometheus_retention_time }} \
{% endif %}
{% if prometheus_retention_size | default('0') != '0' %}
  --storage.tsdb.retention.size={{ prometheus_retention_size }} \
{% endif %}
  --web.console.templates=/etc/prometheus/consoles \
  --web.console.libraries=/etc/prometheus/console_libraries \
{% if prometheus_web_external_url %}
  --web.external-url={{ prometheus_web_external_url }} \
{% endif %}
{% if monitoring_enable_prometheus_backup | default(true) %}
  --web.enable-admin-api \
{% endif %}
  --web.listen-address=0.0.0.0:{{ prometheus_port }}

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
